# image2movie2/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio

# Verzeichnis des Skripts ermitteln
SCRIPT_DIR="$(cd "$(dirname "${BASH_SOURCE[0]}")" && pwd)"

# Setze Ausführungsrechte für create_video.py im gleichen Verzeichnis
chmod +x "$SCRIPT_DIR/create_video.py" || {
    bashio::log.error "Fehler beim Setzen der Ausführungsrechte für create_video.py"
    exit 1
}

# Prüfe, ob Python3 und FFmpeg installiert sind
if ! command -v python3 &> /dev/null; then
    bashio::log.error "Python3 ist nicht installiert. Installiere es zuerst."
    exit 1
elif ! command -v ffmpeg &> /dev/null; then
    bashio::log.error "FFmpeg ist nicht installiert. Installiere es zuerst."
    exit 1
else
    bashio::log.info "Python3 und FFmpeg wurden gefunden."
fi

# Starte das Python-Skript im selben Verzeichnis und logge Ausgaben
bashio::log.info "Starte das Python-Skript $SCRIPT_DIR/create_video.py"
if python3 "$SCRIPT_DIR/create_video.py"; then
    bashio::log.info "Das Python-Skript wurde erfolgreich abgeschlossen."
else
    bashio::log.error "Fehler beim Ausführen des Python-Skripts."
fi