# image2movie2/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio

declare framerate
declare inputpath
declare loglevel
declare revert
declare imageFormat

framerate=$(bashio::config 'framerate')
inputpath=$(bashio::config 'inputpath')
loglevel=$(bashio::config 'loglevel')
revert=$(bashio::config 'revert')
imageFormat=$(bashio::config 'imageFormat')

SCRIPT_DIR="$(cd "$(dirname "${BASH_SOURCE[0]}")" && pwd)"

chmod +x "$SCRIPT_DIR/create_video.py" || {
    bashio::log.error "Fehler beim Setzen der Ausf端hrungsrechte f端r create_video.py"
    exit 1
}

if [ -z "$framerate" ] || [ -z "$inputpath" ] || [ -z "$loglevel" ] || [ -z "$revert" ]; then
    bashio::log.error "Fehlende Konfigurationswerte. Bitte stellen Sie sicher, dass alle Werte konfiguriert sind."
    exit 1
fi

bashio::log.info "Starte das Python-Skript $SCRIPT_DIR/create_video.py mit 端bergebenen Konfigurationswerten"
python3 "$SCRIPT_DIR/create_video.py" "$framerate" "$inputpath" "$loglevel" "$revert" "$imageFormat"

if [ $? -eq 0 ]; then
    bashio::log.info "Das Python-Skript wurde erfolgreich abgeschlossen."
else
    bashio::log.error "Fehler beim Ausf端hren des Python-Skripts."
fi