# image2movie2/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio

declare framerate
declare inputpath
declare loglevel
declare revert

# Konfiguration auslesen
framerate=$(bashio::config 'framerate')
inputpath=$(bashio::config 'inputpath')
loglevel=$(bashio::config 'loglevel')
revert=$(bashio::config 'revert')

# Verzeichnis des Skripts ermitteln
SCRIPT_DIR="$(cd "$(dirname "${BASH_SOURCE[0]}")" && pwd)"

# Setze Ausführungsrechte für create_video.py im gleichen Verzeichnis
chmod +x "$SCRIPT_DIR/create_video.py" || {
    bashio::log.error "Fehler beim Setzen der Ausführungsrechte für create_video.py"
    exit 1
}

# Starte das Python-Skript mit den übergebenen Argumenten und verwende $SCRIPT_DIR
bashio::log.info "Starte das Python-Skript $SCRIPT_DIR/create_video.py mit übergebenen Konfigurationswerten"
python3 "$SCRIPT_DIR/create_video.py" "$framerate" "$inputpath" "$loglevel" "$revert"